# modules/attack_scenarios_module/example_attack_scenarios.json
{
  "database_info": {
    "title": "База знаний КиберНексус",
    "version": "1.1",
    "last_updated": "2025-03-15",
    "description": "База знаний по AI-решениям для кибербезопасности компании КиберНексус с модулем сценариев атак"
  },
  "sections": [
    {
      "id": "attack_scenarios",
      "name": "Сценарии атак",
      "description": "Типовые сценарии кибератак с детальным описанием этапов, техник и мер защиты",
      "subsections": [
        {
          "id": "ransomware_attacks",
          "name": "Атаки программ-вымогателей",
          "description": "Сценарии атак с использованием программ-вымогателей",
          "content": {
            "scenarios": [
              {
                "id": "enterprise_ransomware_attack",
                "name": "Целевая атака программы-вымогателя на предприятие",
                "description": "Сценарий целенаправленной атаки программы-вымогателя на крупное предприятие с проникновением через фишинг, латеральным перемещением и шифрованием критических данных",
                "difficulty_level": "High",
                "impact_level": "Critical",
                "typical_duration": "2-4 недели",
                "detection_complexity": "Medium",
                "mitigation_complexity": "High",
                "tags": ["ransomware", "phishing", "lateral movement", "data encryption"],
                "mitre_attack_ids": ["TA0001", "TA0002", "TA0004", "TA0005", "TA0040"],
                "threat_actors": ["APT группы", "Киберпреступные группировки"],
                "stages": [
                  {
                    "stage_id": "initial_access",
                    "name": "Первоначальный доступ",
                    "description": "Получение первоначального доступа через фишинговое письмо с вредоносным вложением",
                    "order_index": 1,
                    "techniques": [
                      {
                        "name": "Фишинг",
                        "description": "Отправка целевых фишинговых писем сотрудникам с вредоносными макросами в документах Office",
                        "mitre_technique_id": "T1566.001"
                      }
                    ]
                  },
                  {
                    "stage_id": "persistence",
                    "name": "Закрепление",
                    "description": "Установка бэкдора для обеспечения постоянного доступа",
                    "order_index": 2,
                    "techniques": [
                      {
                        "name": "Планировщик задач",
                        "description": "Создание задачи в планировщике для запуска вредоносного кода при загрузке системы",
                        "mitre_technique_id": "T1053.005"
                      }
                    ]
                  },
                  {
                    "stage_id": "lateral_movement",
                    "name": "Латеральное перемещение",
                    "description": "Перемещение по сети для получения доступа к критическим серверам",
                    "order_index": 3,
                    "techniques": [
                      {
                        "name": "Использование учетных данных",
                        "description": "Кража учетных данных и использование протокола WMI для удаленного выполнения команд",
                        "mitre_technique_id": "T1078"
                      }
                    ]
                  },
                  {
                    "stage_id": "privilege_escalation",
                    "name": "Повышение привилегий",
                    "description": "Получение доступа администратора домена",
                    "order_index": 4,
                    "techniques": [
                      {
                        "name": "Kerberoasting",
                        "description": "Атака на службы, использующие Kerberos, для получения хешей учетных записей с повышенными привилегиями",
                        "mitre_technique_id": "T1558.003"
                      }
                    ]
                  },
                  {
                    "stage_id": "exfiltration",
                    "name": "Эксфильтрация данных",
                    "description": "Кража конфиденциальных данных перед шифрованием для двойного вымогательства",
                    "order_index": 5,
                    "techniques": [
                      {
                        "name": "Архивация данных",
                        "description": "Сжатие украденных данных перед эксфильтрацией для минимизации обнаружения",
                        "mitre_technique_id": "T1560"
                      }
                    ]
                  },
                  {
                    "stage_id": "impact",
                    "name": "Воздействие",
                    "description": "Шифрование данных и требование выкупа",
                    "order_index": 6,
                    "techniques": [
                      {
                        "name": "Шифрование данных",
                        "description": "Шифрование критических данных предприятия с использованием сильных алгоритмов шифрования",
                        "mitre_technique_id": "T1486"
                      }
                    ]
                  }
                ],
                "target_assets": [
                  {
                    "asset_type": "Файловые серверы",
                    "asset_description": "Серверы с критическими бизнес-данными, документами и интеллектуальной собственностью"
                  },
                  {
                    "asset_type": "Базы данных",
                    "asset_description": "Серверы баз данных с информацией о клиентах, финансовыми данными и бизнес-информацией"
                  },
                  {
                    "asset_type": "Резервные копии",
                    "asset_description": "Системы резервного копирования, которые атакуются для предотвращения восстановления"
                  }
                ],
                "mitigations": [
                  {
                    "name": "Обучение сотрудников",
                    "description": "Регулярное обучение сотрудников распознаванию фишинговых писем и социальной инженерии",
                    "effectiveness": "Medium",
                    "implementation_complexity": "Low",
                    "stage_id": "initial_access",
                    "product_id": null
                  },
                  {
                    "name": "Сегментация сети",
                    "description": "Разделение сети на изолированные сегменты для предотвращения латерального перемещения",
                    "effectiveness": "High",
                    "implementation_complexity": "Medium",
                    "stage_id": "lateral_movement",
                    "product_id": null
                  },
                  {
                    "name": "ИнтеллектЩит",
                    "description": "Использование продукта ИнтеллектЩит для обнаружения аномального поведения и предотвращения атак программ-вымогателей",
                    "effectiveness": "High",
                    "implementation_complexity": "Medium",
                    "stage_id": null,
                    "product_id": "intellectshield"
                  },
                  {
                    "name": "Изолированные резервные копии",
                    "description": "Хранение резервных копий в изолированной среде, недоступной из основной сети",
                    "effectiveness": "High",
                    "implementation_complexity": "Medium",
                    "stage_id": "impact",
                    "product_id": null
                  }
                ],
                "iocs": [
                  {
                    "ioc_type": "File Hash",
                    "ioc_value": "a1b2c3d4e5f6g7h8i9j0k1l2m3n4o5p6q7r8s9t0",
                    "description": "SHA-256 хеш вредоносного документа с макросом"
                  },
                  {
                    "ioc_type": "Network",
                    "ioc_value": "communication to *.evil-domain.com on port 8080",
                    "description": "Сетевой трафик к командному серверу злоумышленников"
                  },
                  {
                    "ioc_type": "Registry",
                    "ioc_value": "HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
                    "description": "Ключ реестра, используемый для обеспечения персистентности"
                  }
                ],
                "real_world_examples": [
                  {
                    "incident_name": "Атака Colonial Pipeline",
                    "incident_date": "2021-05-07",
                    "affected_organizations": "Colonial Pipeline",
                    "description": "Программа-вымогатель DarkSide атаковала крупнейшего оператора трубопроводов в США, что привело к остановке поставок топлива на восточное побережье США на несколько дней",
                    "outcome": "Компания заплатила выкуп в размере 4.4 миллиона долларов. Позже часть суммы была возвращена ФБР",
                    "lessons_learned": "Необходимость усиленной сегментации сети между ИТ и ОТ системами, улучшение мониторинга безопасности и наличие проверенных планов реагирования на инциденты",
                    "references": "https://www.cisa.gov/sites/default/files/publications/Colonial_Pipeline_Attack_Retrospective_07192021.pdf"
                  }
                ]
              }
            ]
          }
        },
        {
          "id": "apt_attacks",
          "name": "APT-атаки",
          "description": "Сценарии атак продвинутых постоянных угроз (Advanced Persistent Threats)",
          "content": {
            "scenarios": []
          }
        }
      ]
    }
  ]
}
